# roles/paperless-docker-backup/templates/backup-paperless.timer.j2
[Unit]
Description=Run regular backups of paperless

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
