# roles/immich-docker-backup/templates/backup-immich.timer.j2
[Unit]
Description=Run regular backups of immich

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
